# rootfs/etc/services.d/traefik-admin/run
#!/usr/bin/with-contenv bashio

mkdir -p /data/dynamic.d

CERT_RESOLVER="http01"

if [ ! -z "$DNS_PROVIDER" ]; then
    CERT_RESOLVER="dns01"
fi

AUTH_ENDPOINT=$(bashio::config 'authEndpoint')
COOKIE_SECRET=$(bashio::config 'cookieSecret')

if [ ! -z "$AUTH_ENDPOINT" ]; then
    AUTH_ENDPOINT="--AuthEndpoint $AUTH_ENDPOINT"
fi

/web/traefik-admin --ConfigPath /data/dynamic.d --CertResolver $CERT_RESOLVER $AUTH_ENDPOINT --CookieSecret $COOKIE_SECRET