# rootfs/etc/services.d/traefik-admin/run
#!/usr/bin/with-contenv bashio

mkdir -p /data/dynamic.d

HTTPREQ_ENDPOINT=$(bashio::config 'httpreq.url')
HTTPREQ_USERNAME=$(bashio::config 'httpreq.username')
HTTPREQ_PASSWORD=$(bashio::config 'httpreq.password')

CERT_RESOLVER="http01"

if [ ! -z "$HTTPREQ_ENDPOINT" ] && [ ! -z "$HTTPREQ_USERNAME" ] && [ ! -z "$HTTPREQ_PASSWORD" ]; then
    CERT_RESOLVER="dns01"
fi

AUTH_ENDPOINT=$(bashio::config 'authEndpoint')
COOKIE_SECRET=$(bashio::config 'cookieSecret')

if [ ! -z "$AUTH_ENDPOINT" ]; then
    AUTH_ENDPOINT="--AuthEndpoint $AUTH_ENDPOINT"
fi

/web/traefik-admin --ConfigPath /data/dynamic.d --CertResolver $CERT_RESOLVER $AUTH_ENDPOINT --CookieSecret $COOKIE_SECRET