# rootfs/etc/services.d/traefik/run
#!/usr/bin/with-contenv bashio

ACME_EMAIL=$(bashio::config 'email')
mkdir -p /data/dynamic.d

HTTPREQ_ENDPOINT=$(bashio::config 'httpreq.url')
HTTPREQ_USERNAME=$(bashio::config 'httpreq.username')
HTTPREQ_PASSWORD=$(bashio::config 'httpreq.password')

LOGLEVEL=$(bashio::config 'loglevel')

INSECURE_SKIP_VERIFY=$(bashio::config 'insecureSkipVerify')

SKIP_VERIFY=""
if [[ "$INSECURE_SKIP_VERIFY" = "true" ]]; then
    SKIP_VERIFY="--serversTransport.insecureSkipVerify=true"
fi

# Add root CAs
if [  -d /data/cacerts ];then
    rm -r /data/cacerts
fi
mkdir -p /data/cacerts
i=0
CERTCMDLINE="--serversTransport.rootCAs="

for c in $(cat /data/options.json | jq -r '.rootCAs[] | @base64'); do
   echo $c | base64 -d > "/data/cacerts/ca$i.crt"
   CERTCMDLINE="${CERTCMDLINE}/data/cacerts/ca$i.crt,"
   i=$(expr $i + 1)
done

CERTCMDLINE="${CERTCMDLINE%?}"
# If we dont have any cert CAs we empty the variale
if [[ "$CERTCMDLINE" = "--serversTransport.rootCAs" ]]; then
    CERTCMDLINE=
fi

if [ ! -z "$HTTPREQ_ENDPOINT" ] && [ ! -z "$HTTPREQ_USERNAME" ] && [ ! -z "$HTTPREQ_PASSWORD" ]; then
    export HTTPREQ_ENDPOINT="$HTTPREQ_ENDPOINT"
    export HTTPREQ_USERNAME="$HTTPREQ_USERNAME"
    export HTTPREQ_PASSWORD="$HTTPREQ_PASSWORD"
    traefik \
    --log.level=$LOGLEVEL \
    $SKIP_VERIFY \
    --providers.file.directory=/data/dynamic.d \
    --entryPoints.web.address=:80 \
    --entryPoints.websecure.address=:443 \
    --certificatesresolvers.http01.acme.email=$ACME_EMAIL \
    --certificatesresolvers.http01.acme.storage=/data/acme.json \
    --certificatesresolvers.http01.acme.httpchallenge.entrypoint=web \
    --certificatesresolvers.dns01.acme.email=$ACME_EMAIL \
    --certificatesresolvers.dns01.acme.dnschallenge.provider=httpreq \
    --certificatesresolvers.dns01.acme.dnschallenge.delaybeforecheck=1 \
    --certificatesresolvers.dns01.acme.storage=/data/acme.json \
    $CERTCMDLINE
else
    traefik \
    --log.level=$LOGLEVEL \
    $SKIP_VERIFY \
    --providers.file.directory=/data/dynamic.d \
    --entryPoints.web.address=:80 \
    --entryPoints.websecure.address=:443 \
    --certificatesresolvers.http01.acme.email=$ACME_EMAIL \
    --certificatesresolvers.http01.acme.storage=/data/acme.json \
    --certificatesresolvers.http01.acme.httpchallenge.entrypoint=web \
    $CERTCMDLINE
fi


